# templates/systemd.service.erb
[Unit]
Description=<%=@name%>
After=docker.service
Requires=docker.service
 
[Service]
TimeoutStartSec=<%=@timeout%>
ExecStartPre=/usr/bin/docker pull <%=@image%>:<%=@tag%>
ExecStartPre=/usr/bin/docker create --rm --name <%=@name%> <%=@options.map{|k,v| "--#{k} #{v}"}.join(' ')%> <%=@image%>:<%=@tag%>

ExecStart=/usr/bin/docker start <%=@name%>

ExecStop=/bin/docker stop <%=@name%>
ExecStopPost=/bin/docker rm -f <%=@name%>

ExecReload=/bin/docker restart <%=@name%>

[Install]
WantedBy=multi-user.target
